# debian/postinst
#!/bin/sh

#DEBHELPER#

rc-update del gps-nokia-n900 default || true
